# config/display2/etc/cron.d/easycast-refresh-display
*/30 * * * * pi /home/pi/easycast/bin/refresh-displays
