# config/displays-lxde/etc/cron.d/easycast-refresh-display.tpl
*/30 * * * * {{{station_config.easycast_user}}} {{{station_config.easycast_user_home}}}/easycast/bin/refresh-displays
